# docker/conf/etc/services.d/speedtest/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc php /config/www/artisan queue:work --timeout=120 >> /config/log/speedtest/queue.log